# systemd/verdify-forecast-page.service
[Unit]
Description=Regenerate generated lab.verdify.ai content after forecast refresh
After=verdify-timescaledb.service

[Service]
Type=oneshot
Environment=TZ=America/Denver
Environment=PYTHON=/srv/greenhouse/.venv/bin/python
ExecStart=/srv/verdify/scripts/publish-site-content.sh --reason forecast
StandardOutput=append:/var/local/verdify/state/forecast-page.log
StandardError=append:/var/local/verdify/state/forecast-page.log
User=jason
Nice=10
