# systemd/verdify-site-poll.timer
[Unit]
Description=Verdify site — 10-second vault polling timer

[Timer]
OnBootSec=30s
OnUnitActiveSec=10s
AccuracySec=1s
Unit=verdify-site-poll.service

[Install]
WantedBy=timers.target
