# systemd/verdify-plan-publish.service
[Unit]
Description=Publish the latest planner output to verdify.ai
After=verdify-timescaledb.service docker.service

[Service]
Type=oneshot
# Use MDT date semantics — plan_id timestamps are in America/Denver, so the
# plan page we publish is "today in Denver" not "today in UTC".
Environment=TZ=America/Denver
Environment=PYTHON=/srv/greenhouse/.venv/bin/python
ExecStart=/srv/verdify/scripts/publish-site-content.sh --reason planner
StandardOutput=journal
StandardError=journal
User=jason
Nice=10
