# systemd/verdify-setpoint-server.service
[Unit]
Description=Verdify Setpoint Server (ESP32 HTTP endpoint + grow light control)
After=network.target docker.service
Wants=network.target

[Service]
Type=simple
User=jason
ExecStart=/srv/greenhouse/.venv/bin/python3 /srv/verdify/scripts/setpoint-server.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
