# systemd/verdify-mcp.service
[Unit]
Description=Verdify MCP Server (Greenhouse tools for Agent Iris)
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=jason
WorkingDirectory=/srv/verdify
ExecStart=/srv/greenhouse/.venv/bin/python mcp/server.py
Restart=always
RestartSec=10
StandardOutput=append:/srv/verdify/state/mcp-server.log
StandardError=append:/srv/verdify/state/mcp-server.log

[Install]
WantedBy=multi-user.target
