# systemd/verdify-forecast.service
[Unit]
Description=Verdify Forecast Sync (Open-Meteo 72h)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
User=jason
ExecStart=/srv/verdify/scripts/forecast-sync.py
StandardOutput=append:/srv/verdify/state/forecast-sync.log
StandardError=append:/srv/verdify/state/forecast-sync.log

[Install]
WantedBy=multi-user.target
