# systemd/verdify-ingestor.service
[Unit]
Description=Verdify ESP32 data ingestor
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=jason
WorkingDirectory=/srv/verdify/ingestor
KillMode=control-group
ExecStartPre=/bin/bash -c 'pgrep -f "^/srv/greenhouse/.venv/bin/python ingestor.py" | xargs -r kill; sleep 1; true'
ExecStart=/srv/greenhouse/.venv/bin/python ingestor.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
